# home/.config/autostart/alacritty.desktop
[Desktop Entry]
Exec=alacritty --class alacritty-tmux -e tmux-workspace.sh
Icon=
Name=alacritty
Path=
Terminal=False
Type=Application
