# home/.config/autostart/clipit.desktop
[Desktop Entry]
Exec=clipit -n
Name=clipit
Type=Application
Version=1.0
