# home/.config/autostart/krunner.desktop
[Desktop Entry]
Exec=/usr/bin/krunner -d
Icon=
Name=krunner
Path=
Terminal=False
Type=Application
