# kde/.config/autostart/foot.desktop
[Desktop Entry]
Exec=foot -a foot-tmux -e tmux-workspace.sh
Icon=
Name=foot
Path=
Terminal=False
Type=Application
