# kde/.config/autostart/PCManFM-Qt File Manager.desktop
[Desktop Entry]
Exec=pcmanfm-qt -d %U
Icon=pcmanfm-qt
Name=PCManFM-Qt File Manager
Path=
Terminal=False
Type=Application
