# kde/.config/autostart/run-clipcatd.sh.desktop
[Desktop Entry]
Exec=~/bin/run-clipcatd.sh
Icon=
Name=run-clipcatd.sh
Path=
Terminal=False
Type=Application
